# update-reactive.timer
[Unit]
Description=Monthly update of reactive firewall rules

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target

